# build files/media-backup.service

[Unit]
Description=Media backup & snapshot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/media-backup.sh
